# package/plugin-gargoyle-usb-printer/files/etc/hotplug.d/usb/80-usb-printer
#!/bin/sh

if [ "$ACTION" = "add" ] ||  [ "$ACTION" = "remove" ] ; then
	sh /usr/lib/gargoyle/configure_printer.sh
fi

