# package/plugin-gargoyle-i18n/files/www/etc/uci-defaults/zzz-plugin-gargoyle-i18n
#!/bin/sh

. /usr/lib/gargoyle/i18nServices.sh
change_menu_language $(uci get gargoyle.global.fallback_lang)
restart_lang_services
