# package/plugin-gargoyle-wifi-survey/files/etc/uci-defaults/plugin-gargoyle-wifi-survey
#!/bin/sh

uci set gargoyle.display.system_wifi_survey='WiFi Survey'
uci set gargoyle.scripts.system_wifi_survey='wifi_survey.sh'
uci set gargoyle.system.wifi_survey='380'
uci commit
