# package/plugin-gargoyle-cron/files/etc/uci-defaults/plugin-gargoyle-cron
#!/bin/sh

uci set gargoyle.display.system_cron='Cron jobs'
uci set gargoyle.scripts.system_cron='cron.sh'
uci set gargoyle.system.cron='330'
uci commit
