# package/plugin-gargoyle-wifi-schedule/files/etc/uci-defaults/plugin-gargoyle-wifi-schedule
#!/bin/sh

uci set gargoyle.display.system_wifi_schedule='WiFi Schedule'
uci set gargoyle.scripts.system_wifi_schedule='wifi_schedule.sh'
uci set gargoyle.system.wifi_schedule='370'
uci commit
