# package/plugin-gargoyle-usb-storage/files/etc/hotplug.d/block/80-usb-storage
#!/bin/sh

case $DEVICENAME in
	mtdblock*);;
	*)
		/etc/init.d/usb_storage restart
		;;
esac
