# package/ddns-gargoyle/files/etc/hotplug.d/iface/25-ddns
#!/bin/sh

if [ "$INTERFACE" = "wan" ] && [ "$ACTION" = "ifup" ] ; then
	sleep 5
	/usr/bin/ddns_gargoyle -P /etc/ddns_providers.conf -C /etc/ddns_gargoyle.conf
fi
