# package/plugin-gargoyle-usb-storage/files/etc/hotplug.d/usb/15-usb-storage
#!/bin/sh

if [ "$ACTION" = "add" ] || [ "$ACTION" = "remove" ] ; then
	if [ "$ACTION" = "add" ] ; then sleep 3; fi
	/etc/init.d/usb_storage restart
fi
