# package/portmap/files/portmap.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2008 OpenWrt.org

START=19
STOP=19

NAME=portmap
PROG=/usr/sbin/$NAME
DEFAULT=/etc/default/$NAME

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	start-stop-daemon -q -S -x $PROG -- $OPTIONS
}

stop() {
	start-stop-daemon -q -K -x $PROG -- $OPTIONS
}
