# package/plugin-gargoyle-email-notifications/files/etc/uci-defaults/plugin-gargoyle-email-notifications
#!/bin/sh

uci set gargoyle.display.email='Email'
uci set gargoyle.scripts.email='email.sh'
uci set gargoyle.system.email='385'
uci commit
