# package/plugin-gargoyle-email-notifications/files/etc/uci-defaults/plugin-gargoyle-email-notifications
#!/bin/sh

uci set gargoyle.display.system_email_notifications='Email'
uci set gargoyle.scripts.system_email_notifications='email.sh'
uci set gargoyle.system.email_notifications='380'
uci commit
