# package/plugin-gargoyle-openvpn/files/etc/hotplug.d/iface/27-openvpn
#!/bin/sh

config_load "network"
config_get wan_proto "wan" "proto"

if [ "$INTERFACE" = "wan" ] && [ "$ACTION" = "ifup" ] 
then
	sleep 15
	/etc/init.d/openvpn restart
fi


