# package/plugin-gargoyle-minidlna/files/etc/uci-defaults/plugin-gargoyle-minidlna
#!/bin/sh

uci set gargoyle.display.system_minidlna='DLNA'
uci set gargoyle.scripts.system_minidlna='minidlna.sh'
uci set gargoyle.system.minidlna='352'
uci commit gargoyle
