# package/gargoyle/files/etc/hotplug.d/net/20-client_bridge
if [ "$ACTION" = "add" -o "$ACTION" = "register" ]; then
	is_ath=$(echo $INTERFACE | grep ath)
	if [ -n "$is_ath" ] ; then
		managed=$(iwconfig $INTERFACE 2>/dev/null | grep "Mode:Managed")
		if [ -n "$managed" ] ; then
			sta_section=$(uci show wireless | grep "mode=sta$" | sed 's/\.mode=sta$//g' )			
			if [ -n "$sta_section" ] ; then
				client_bridge=$(uci get "$sta_section.client_bridge")
				if [ "$client_bridge" = "1" ] ; then
					ebtables -t nat -F PREROUTING
					ebtables -t nat -F POSTROUTING
					ebtables -t broute -F BROUTING
					
					ebtables -t nat -A PREROUTING  --in-interface  $INTERFACE -j arpnat --arpnat-target ACCEPT
					ebtables -t nat -A POSTROUTING --out-interface $INTERFACE -j arpnat --arpnat-target ACCEPT
					ebtables -t broute -A BROUTING --in-interface  $INTERFACE --protocol 0x888e -j DROP
				fi
			fi
		fi
	fi
fi
