# package/gargoyle/files/etc/hotplug.d/net/20-client_bridge
is_ath=$(echo $INTERFACE | grep ath)
if [ -n "$is_ath" ] ; then
	if [ "$ACTION" = "add" -o "$ACTION" = "register" ] ; then
		sleep 2
		managed=$(iwconfig $INTERFACE 2>/dev/null | grep "Mode:Managed")
		if [ -n "$managed" ] ; then
			sta_section=$(uci show wireless | grep "mode=sta$" | sed 's/\.mode=sta$//g' )			
			if [ -n "$sta_section" ] ; then
				client_bridge=$(uci get "$sta_section.client_bridge")
				if [ "$client_bridge" = "1" ] ; then
					ebtables -t nat -F PREROUTING
					ebtables -t nat -F POSTROUTING
					ebtables -t broute -F BROUTING
				
					bridge_ip=$(uci get network.lan.ipaddr)
					if [ -n "$bridge_ip" ] ; then
						#always arpnat DHCP
						ebtables -t nat -A PREROUTING  --in-interface  $INTERFACE  -p IPv4 --ip-protocol UDP --ip-destination-port 67:68 -j arpnat --arpnat-target ACCEPT
						ebtables -t nat -A PREROUTING  --in-interface  $INTERFACE  -p IPv4 --ip-protocol UDP --ip-source-port      67:68 -j arpnat --arpnat-target ACCEPT
						
						#but in order to allow access to router from external network don't arpnat anything else destined for bridge device
						ebtables -t nat -A PREROUTING  -p IPv4 --ip-dst $bridge_ip -j redirect --redirect-target ACCEPT
					fi
					
					ebtables -t nat -A PREROUTING  --in-interface  $INTERFACE -j arpnat --arpnat-target ACCEPT
					ebtables -t nat -A POSTROUTING --out-interface $INTERFACE -j arpnat --arpnat-target ACCEPT
					ebtables -t broute -A BROUTING --in-interface  $INTERFACE --protocol 0x888e -j DROP
					
					mkdir -p /tmp/client_bridge
					touch "/tmp/client_bridge/$INTERFACE"
				fi
			fi
		fi
	fi
	if [ "$ACTION" = "remove" -o "$ACTION" = "register" ] ; then
		if [ -e "/tmp/client_bridge/$INTERFACE" ] ; then
			ebtables -t nat -F PREROUTING
			ebtables -t nat -F POSTROUTING
			ebtables -t broute -F BROUTING
		fi			
	fi
fi
