# package/gargoyle/files/etc/hotplug.d/iface/99-qmincm
#!/bin/sh

[ $ACTION = "ifup" ] || exit 0
[ $INTERFACE = "wan" ] || exit 0

case $(uci -q get network.wan.proto) in
	qmi|ncm|mbim)
	uci -P /var/state set network.wan.device=$DEVICE
	;;
esac
exit 0
