# package/plugin-gargoyle-usb-printer/files/etc/uci-defaults/init_gargoyle_printer_page
#!/bin/sh

uci set gargoyle.display.system_printers="USB Printer"
uci set gargoyle.scripts.system_printers="printers.sh"
uci set gargoyle.system.printers="360"
uci commit
