# package/plugin-gargoyle-openvpn/files/etc/uci-defaults/init_gargoyle_openvpn_plugin
#!/bin/sh

uci set gargoyle.display.connection_openvpn="OpenVPN"
uci set gargoyle.scripts.connection_openvpn="openvpn.sh"
uci set gargoyle.connection.openvpn="225"
uci commit
