# package/plugin-gargoyle-initd/files/etc/uci-defaults/plugin-gargoyle-initd
#!/bin/sh

uci set gargoyle.display.system_initd='Services'
uci set gargoyle.scripts.system_initd='initd.sh'
uci set gargoyle.system.initd='310'
uci commit
