# package/plugin-gargoyle-adblock/files/etc/uci-defaults/plugin-gargoyle-adblock
#!/bin/sh

uci set gargoyle.display.system_adblock='Adblock'
uci set gargoyle.scripts.system_adblock='adblock.sh'
uci set gargoyle.system_adblock='357'
uci commit
