# package/plugin-gargoyle-adblock/files/etc/uci-defaults/plugin-gargoyle-adblock
#!/bin/sh

uci set gargoyle.display.firewall_adblock='Adblock'
uci set gargoyle.scripts.firewall_adblock='ablock.sh'
uci set gargoyle.firewall.adblock='120'
uci commit
