# package/plugin-gargoyle-webshell/files/etc/uci-defaults/plugin-gargoyle-webshell
#!/bin/sh

uci set gargoyle.display.system_webshell='Webshell'
uci set gargoyle.scripts.system_webshell='webshell.sh'
uci set gargoyle.system.webshell='312'
uci commit
