# package/plugin-gargoyle-usb-printer/files/etc/hotplug.d/iface/80-usb-printer
#!/bin/sh
if [ "$INTERFACE" = "lan" ] && [ "$ACTION" = "ifup" ] ; then
	sh /usr/lib/gargoyle/configure_printer.sh
fi

