# package/plugin-gargoyle-spectrum-analyser/files/etc/uci-defaults/plugin-gargoyle-spectrum-analyser
#!/bin/sh

uci set gargoyle.display.system_spectrum_analyser='Spectrum Analyser'
uci set gargoyle.scripts.system_spectrum_analyser='spectrum_analyser.sh'
uci set gargoyle.system.spectrum_analyser='355'
uci commit
